# pifancontrol.service
[Unit]
Description=Service for pifancontrol

[Service]
ExecStart=/usr/local/sbin/
Restart=on-failure

[Install]
WantedBy=default.target
